# cleanup.ini
[program:cleanup]
priority = 100
command = bash -c "while i=1; do sleep 86400; cat /dev/null > /njmon/hmc_collector.log && cat /dev/null >/ njmon/njmond.log; done"
startsecs = 0
autorestart = true
