# influxdb.ini
[program:influxdb]
command=/usr/bin/influxd -pidfile /var/run/influxdb/influxd.pid -config /etc/influxdb/influxdb.conf
process_name = influxd
autostart  = true
autorestart = true
