# njmon_collector.ini
[program:njmon_collector]
command=/njmon/njmond.py /njmon/njmond.conf
process_name = njmond.py
autostart  = true
autorestart = true
