# nmon2influxdb.ini
[program:nmon2influxdb]
priority = 100
command = bash -c "while i=1; do sleep 1200; /njmon/hmc_collector.sh >> /njmon/hmc_collector.log; done"
startsecs = 0
autorestart = true
